# debian/dwww.cron.weekly
#!/bin/sh
# $Id: dwww.cron.weekly,v 1.2 2006-05-30 18:56:45 robert Exp $
#
# This script creates an index of Debian documentation registered
# with doc-base. See dwww-index++(8) man page for details.

# check if swish++ is installed
test -x /usr/bin/index++ || exit 0

# check if dwww is still installed
test -x /usr/sbin/dwww-index++ || exit 0

set -e

dwww-index++ > /dev/null

exit 0
